# examples/price_estimation/DockerFiles/tendermint/config_service
#!/sbin/openrc-run


description="Executes an abitrary python script."


start() {
        python3 /tendermint/config_loader.py 
}